# etc/cron.daily/dirvish-cronjob
#!/bin/sh
#
# dirvish's cron script
#
# This script runs the daily dirvish cronjob.

test -x "$(which dirvish-cronjob)" || exit 0
dirvish-cronjob
